# config/systemd/user/xautolock@.service
[Service]
Environment=DISPLAY=%I
ExecStart=/usr/bin/xautolock -time 5 -locker xtrlock
Restart=on-failure

[Install]
WantedBy=graphics.target
